# systemd/audio-virtual-autoset.service
[Unit]
Description=Watch Bluetooth mode and auto-link virtual audio
After=sound.target pipewire.service

[Service]
Type=simple
ExecStart=%h/source/audio-virtual-linker/bin/virtual-autoset.sh
Restart=on-failure
RestartSec=2
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
